# discord-bot/com.english-teacher.discord-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.english-teacher.discord-bot</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/yuki/.bun/bin/bun</string>
    <string>run</string>
    <string>discord-bot/src/index.ts</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/yuki/Work/ai-agents/english-teacher-agent</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/yuki</string>
    <key>PATH</key>
    <string>/Users/yuki/.local/bin:/Users/yuki/.bun/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/yuki/Work/ai-agents/english-teacher-agent/discord-bot/logs/stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/yuki/Work/ai-agents/english-teacher-agent/discord-bot/logs/stderr.log</string>
</dict>
</plist>
